# yoga_fan.service
[Unit]
Description=Yoga Fan Control
[Service]
Type=forking
ExecStart=/usr/bin/yoga_fan
[Install]
WantedBy=multi-user.target